# deploy/systemd/agent-ops-weekly-report.service
[Unit]
Description=Agent Ops Weekly Report Generator
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/angad-option-desk-v1
EnvironmentFile=-/etc/angad-option-desk.env
ExecStart=/home/ubuntu/angad-option-desk-v1/venv/bin/python -m app.agent_ops.cli weekly-report

[Install]
WantedBy=multi-user.target
