# deploy/systemd/research-weekly.service
[Unit]
Description=Research Intelligence Weekly Strategy Review
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/angad-option-desk-v1
EnvironmentFile=-/etc/angad-option-desk.env
ExecStart=/home/ubuntu/angad-option-desk-v1/venv/bin/python -m app.research_intelligence.cli weekly-strategy-review

[Install]
WantedBy=multi-user.target
