# deploy/systemd/telegram-control-room.service
[Unit]
Description=Telegram Control Room Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/angad-option-desk-v1
EnvironmentFile=-/etc/angad-option-desk.env
ExecStart=/home/ubuntu/angad-option-desk-v1/venv/bin/python -m app.telegram_control_room.cli run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
